# server/package/DEBIAN/postinst
#!/bin/bash
set -e

# Verify Node.js is installed and meets minimum version
if ! command -v node &>/dev/null; then
  echo "ERROR: nodejs not found. Install it first:"
  echo "  curl -fsSL https://deb.nodesource.com/setup_20.x | sudo bash -"
  echo "  sudo apt install -y nodejs"
  exit 1
fi
NODE_MAJOR=$(node -e "process.stdout.write(process.versions.node.split('.')[0])")
if [ "$NODE_MAJOR" -lt 18 ]; then
  echo "ERROR: Node.js >= 18 required (found $(node -v))"
  exit 1
fi

# Create system user if absent
if ! id thor &>/dev/null; then
  useradd -r -s /bin/false -d /etc/thor-server thor
fi

# Create data directory and seed rooms.json only on first install
mkdir -p /etc/thor-server/data
if [ ! -f /etc/thor-server/data/rooms.json ]; then
  cp /usr/share/thor-server/rooms.json.default /etc/thor-server/data/rooms.json
fi
chown -R thor:thor /etc/thor-server/data
chmod 750 /etc/thor-server/data

# Create log directory
mkdir -p /var/log/thor-server
chown thor:thor /var/log/thor-server
chmod 750 /var/log/thor-server

# Server is a single self-contained bundle — nothing to install at runtime.
chown -R root:thor /usr/lib/thor-server
chmod -R 755 /usr/lib/thor-server

# Generate .env on first install only (conffiles preserves it on upgrade)
if [ ! -f /etc/thor-server/.env ]; then
  cat > /etc/thor-server/.env << EOF
NODE_ENV=production
PORT=3001
JWT_SECRET=$(openssl rand -hex 32)
ALLOWED_ORIGINS=https://assets.devshram.com
FRONTEND_BASE_URL=https://assets.devshram.com/projects/thor
DATA_DIR=/etc/thor-server/data
LOG_DIR=/var/log/thor-server
EOF
  chmod 600 /etc/thor-server/.env
  chown thor:thor /etc/thor-server/.env
  echo ""
  echo "================================================"
  echo "  Thor installed!"
  echo "  Sign in to the frontend with:"
  echo "    username: root"
  echo "    password: thor"
  echo "  Change the password immediately."
  echo "================================================"
  echo ""
fi

# Upgrade path: ensure JWT_SECRET exists in pre-existing .env files
if [ -f /etc/thor-server/.env ] && ! grep -q '^JWT_SECRET=' /etc/thor-server/.env; then
  echo "JWT_SECRET=$(openssl rand -hex 32)" >> /etc/thor-server/.env
fi

# Seed empty users.json on first install
if [ ! -f /etc/thor-server/data/users.json ]; then
  echo '{"users":[]}' > /etc/thor-server/data/users.json
  chown thor:thor /etc/thor-server/data/users.json
  chmod 640 /etc/thor-server/data/users.json
fi

# Stop any stray thor-server process that's not under systemd (e.g. a manual
# `npm run dev` or leftover from a previous install). Without this, the new
# systemd unit can't bind to PORT and would silently fail to start.
PORT=$(grep -E '^PORT=' /etc/thor-server/.env | cut -d= -f2)
PORT="${PORT:-3001}"
if command -v lsof &>/dev/null; then
  STRAY_PIDS=$(lsof -tiTCP:"$PORT" -sTCP:LISTEN 2>/dev/null | grep -v "^$$\$" || true)
  if [ -n "$STRAY_PIDS" ]; then
    echo "Found process(es) listening on port $PORT — stopping them"
    # shellcheck disable=SC2086
    kill $STRAY_PIDS 2>/dev/null || true
    sleep 1
    # Force-kill anything still alive
    # shellcheck disable=SC2086
    kill -9 $STRAY_PIDS 2>/dev/null || true
  fi
elif command -v fuser &>/dev/null; then
  fuser -k "${PORT}/tcp" 2>/dev/null || true
fi

# Enable and start services. Failures here must not abort the install — the
# package is on disk and a service issue can be debugged with journalctl.
if [ -d /run/systemd/system ]; then
  systemctl daemon-reload                              || true
  systemctl enable thor-server thor-update.timer       || true
  systemctl restart thor-server                        || echo "WARN: thor-server failed to start. Check: journalctl -u thor-server"
  systemctl start thor-update.timer                    || true
fi
