# server/package/lib/systemd/system/thor-update.timer
[Unit]
Description=Thor Server - Auto-update poll timer
Wants=thor-update.service

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Unit=thor-update.service

[Install]
WantedBy=timers.target
